# assignment04/90-keyboard-mod-remove.rules
# AUTHOR: MYPARK
ACTION=="remove", ENV{ID_MODEL}=="QEMU_USB_Keyboard",  RUN+="/sbin/modprobe -r keyboard_mod"

